# deploy/gunicorn.service
[Unit]
Description=LyomasPharmacy (tweb) Gunicorn
After=network.target

[Service]
User=softwarehouse
Group=www-data
WorkingDirectory=/var/www/tweb
EnvironmentFile=/var/www/tweb/.env
ExecStart=/var/www/tweb/venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:8001 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/tweb/access.log \
    --error-logfile /var/log/tweb/error.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
